# src/cjn_PiTools/deployment_files/PiOLED_server.service
[Unit]
Description=PiOLED display service

[Service]
ExecStart=/path_to/PiOLED --service
User=me
Group=me

[Install]
WantedBy=multi-user.target
